# etc/supervisord.conf.in
[unix_http_server]
file=/tmp/supervisor_comet.sock ; path to your socket file

;[inet_http_server]
;port = 127.0.0.1:9001

[supervisorctl]
serverurl=unix:///tmp/supervisor_comet.sock ; use a unix:// URL for a unix socket
;[supervisorctl]
;serverurl = http://localhost:9001


[supervisord]
logfile = ${buildout:directory}/log/supervisor/supervisord.log
logfile_maxbytes = 50MB
logfile_backups = 10
loglevel = info
pidfile = ${buildout:directory}/tmp/supervisord.pid
nodaemon = false


[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

#[include]
#files = *.supervisor

[program:web]
command=${buildout:directory}/bin/python ${buildout:directory}/src/app.py
stderr_logfile = ${buildout:directory}/log/supervisor/web-stderr.log
stdout_logfile = ${buildout:directory}/log/supervisor/web-stdout.log

[program:online_user]
command=${buildout:directory}/bin/python ${buildout:directory}/src/subscribers/online_user.py
stderr_logfile = ${buildout:directory}/log/supervisor/online_user-stderr.log
stdout_logfile = ${buildout:directory}/log/supervisor/online_user-stdout.log

[program:room_online_user]
command=${buildout:directory}/bin/python ${buildout:directory}/src/subscribers/room_online_user.py
stderr_logfile = ${buildout:directory}/log/supervisor/room_online_user-stderr.log
stdout_logfile = ${buildout:directory}/log/supervisor/room_online_user-stdout.log

[program:room]
command=${buildout:directory}/bin/python ${buildout:directory}/src/subscribers/room.py
stderr_logfile = ${buildout:directory}/log/supervisor/room-stderr.log
stdout_logfile = ${buildout:directory}/log/supervisor/room-stdout.log

[program:gc]
command=${buildout:directory}/bin/python ${buildout:directory}/src/gc.py
stderr_logfile = ${buildout:directory}/log/supervisor/gc-stderr.log
stdout_logfile = ${buildout:directory}/log/supervisor/gc-stdout.log
